# systemd/vibedocs.service
[Unit]
Description=VibeDocs - Documentation Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# CONFIGURE: Set to your VibeDocs installation directory
WorkingDirectory=/path/to/vibedocs
# CONFIGURE: Set to your tsx binary (run: which tsx or npx which tsx)
ExecStart=/path/to/vibedocs/node_modules/.bin/tsx src/server.ts
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=10
Environment=NODE_ENV=production
# CONFIGURE: Set to the directory containing your project folders
Environment=VIBEDOCS_ROOT=/path/to/your/projects
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vibedocs

[Install]
WantedBy=default.target
